# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[program:app]
command=uvicorn src.main:app --host 0.0.0.0 --port 8000 --workers 4
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log

[program:worker]
command=celery -A src.tasks.celery_app worker --loglevel=info --concurrency=4
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log

[program:beat]
command=celery -A src.tasks.celery_app beat --loglevel=info
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/beat.err.log
stdout_logfile=/var/log/beat.out.log
